# .config/systemd/user/hyprsunset-auto.service
[Unit]
Description=Auto-toggle hyprsunset at 7 PM

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/hyprsunsetctl auto
